# web/profexor/settings/supervisord-db.conf
[program:createdb]
command=/createdb.sh
numprocs=1
autostart=true
